# supervisor/example.conf
[program:myapp]
directory=/home/example/app/myapp
command=/usr/bin/node myapp.js
user=app
